# packaging/rpm/sakms-node.service
[Unit]
Description=sakms worker node daemon
Documentation=https://github.com/curtiswtaylorjr/sakms
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Security-hardening addendum: runs as a dedicated, non-root system user
# (created by the RPM's %pre scriptlet) instead of root -- a compromised or
# buggy daemon then has only this user's own filesystem permissions, not
# unbounded root access. Read access to actual media is granted separately
# (e.g. supplementary group membership matching the operator's own media
# ownership scheme -- see the plan's own execution-model discussion), which
# is a per-host operator step, not something this unit file can generalize.
User=sakms-node
Group=sakms-node
ExecStart=/usr/bin/sakms-node -config /etc/sakms-node/config.json
Restart=on-failure
RestartSec=10
ProtectSystem=full
ReadWritePaths=/etc/sakms-node

# CPU governor cgroup delegation (Option C, see cmd/sakms-node/resourcegov.go
# and the node-resource-governor plan's Stage 3 Packaging bullet): without
# Delegate=yes, systemd keeps this unit's cgroup subtree root-owned, so the
# non-root sakms-node user above has no permission to create/manage its own
# leaf cgroup or write cpu.max -- capabilityDetect() would report
# "unavailable" and the cap could never actually enforce. Delegate=yes hands
# ownership of this unit's cgroup subtree to sakms-node, which an empirical
# spike (see CG-S1) already proved is sufficient for the daemon to create its
# leaf cgroup and enforce the CPU cap without any additional privilege.
Delegate=yes

# mediaRoots local control socket (sakms-node-tray UI, see the ralplan
# mediaRoots-UI plan's socket-perms/Option-A section): the daemon listens on
# a unix socket at /run/sakms-node/control.sock. RuntimeDirectory gives it a
# tmpfs dir that systemd creates before start and removes on graceful stop;
# RuntimeDirectoryMode=0750 (owner sakms-node) is deliberately NOT group-
# traversable by sakms-media-config yet on its own -- the daemon chgrp's the
# dir (and chmod/chgrp's the socket file to 0660 sakms-media-config) itself
# at startup via os.Chown/os.Chmod (Option A), immediately after net.Listen
# succeeds, which is the fail-closed startup window the plan documents.
# SupplementaryGroups=sakms-media-config makes the daemon a member of that
# shared group (the daemon's *primary* group stays sakms-node -- config.json
# stays 0600 sakms-node:sakms-node, untouched by this), which POSIX chown(2)
# requires for the daemon's own unprivileged chgrp-to-group calls above to
# succeed. Group membership in sakms-media-config is the sole authorization
# boundary for the control socket (see sakms-node.spec's %post and
# sakms-node.sysusers.conf for how the shared group is created and how the
# console/desktop user is added to it).
RuntimeDirectory=sakms-node
RuntimeDirectoryMode=0750
SupplementaryGroups=sakms-media-config

[Install]
WantedBy=multi-user.target
